# workflow-timer.service
[Unit]
Description=Workflow Timer

[Service]
ExecStart=/usr/bin/workflow-timer

[Install]
WantedBy=multi-user.target
